# host-agent/mercer-host-agent.service
[Unit]
Description=Mercer Host Agent — sidecar process manager
After=network.target

[Service]
Type=simple
User=%i
# --- Замените пути под вашу систему ---
WorkingDirectory=/path/to/mercer/host-agent
ExecStart=/path/to/mercer/host-agent/.venv/bin/python agent.py
Restart=on-failure
RestartSec=5

# Переменные окружения
Environment=HOST_AGENT_PORT=9090
Environment=HOST_AGENT_TOKEN=changeme
Environment=SIDECAR_DIR=/path/to/mercer/pdf-sidecar

# Безопасность: агент работает только на localhost, процесс не privileged
# Раскомментируйте если используете systemd 235+
# PrivateTmp=true
# NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
